# contrib/systemd/claudegate.service
[Unit]
Description=Claudegate - Anthropic API to AWS Bedrock proxy
After=network.target

[Service]
Type=simple
User=%i
ExecStart=%h/.local/bin/claudegate
Restart=on-failure
RestartSec=5
Environment=AWS_REGION=us-west-2
Environment=CLAUDEGATE_LOG_LEVEL=INFO

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
