# contrib/systemd/claudegate.service
[Unit]
Description=Claudegate - API proxy for AWS Bedrock and GitHub Copilot
After=network.target

[Service]
Type=simple
User=%i
ExecStart=%h/.local/bin/claudegate
Restart=on-failure
RestartSec=5
# Add AWS_REGION here if using the Bedrock backend
# Environment=AWS_REGION=us-west-2
Environment=CLAUDEGATE_LOG_LEVEL=INFO
Environment=NO_COLOR=1

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
